# conf/generic/etc/rc.local
#!/bin/sh -e
#
# rc.local
#
# This script is executed at the end of each multiuser runlevel.
# Make sure that the script will "exit 0" on success or any other
# value on error.

if [ -f /home/spr/.spr-setup-done ]; then
  /home/spr/run.sh
else
  /home/spr/setup.sh
  /home/spr/run.sh
fi

exit 0
